# schedules/templates/com.adjutant.channel-session.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.adjutant.channel-session</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd ~/adjutant &amp;&amp; claude --channels plugin:telegram@claude-plugins-official --model opus --effort medium --remote-control "Adjutant"</string>
    </array>
    <key>WorkingDirectory</key>
    <string>$HOME/adjutant</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string></string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>$HOME/adjutant/logs/channel-session.log</string>
    <key>StandardErrorPath</key>
    <string>$HOME/adjutant/logs/channel-session.err</string>
</dict>
</plist>
